# root/etc/s6-overlay/s6-rc.d/atlantis/run
#!/command/execlineb -P
/usr/local/bin/atlantis server --config /config/etc/server.yaml