# root/etc/services.d/atlantis/run
#!/command/execlineb -P

s6-setuidgid abc /usr/local/bin/atlantis server --config /config/etc/server.yaml
